# templates/docker_exporter.service.j2
[Unit]
Description=Docker Exporter

[Service]
ExecStart=/opt/docker_exporter/docker_status_exporter_amd64 $OPTIONS

[Install]
WantedBy=multi-user.target
